# BMO-setup/docker/bmo-backup.service
[Unit]
Description=BMO Daily Backup to Google Drive
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=patrick
Group=patrick
WorkingDirectory=/home/patrick/bmo
ExecStart=/bin/bash /home/patrick/bmo/backup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
