# BMO-setup/docker/bmo.service
[Unit]
Description=BMO AI Assistant
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=patrick
Group=patrick
WorkingDirectory=/home/patrick/bmo
Environment=PATH=/home/patrick/bmo/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
EnvironmentFile=/home/patrick/bmo/.env
ExecStart=/home/patrick/bmo/venv/bin/python app.py
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

# Hardware access
SupplementaryGroups=video audio i2c gpio spi input

[Install]
WantedBy=multi-user.target
