# BMO-setup/pi/kiosk/bmo-social-bot.service
[Unit]
Description=BMO Social Discord Bot — Music, Chat & Fun
After=network-online.target bmo.service
Wants=network-online.target

[Service]
Type=simple
User=patrick
WorkingDirectory=/home/patrick/bmo
EnvironmentFile=/home/patrick/bmo/.env
ExecStart=/home/patrick/bmo/venv/bin/python discord_social_bot.py
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=append:/home/patrick/bmo/data/logs/social-bot.log
StandardError=append:/home/patrick/bmo/data/logs/social-bot.log

[Install]
WantedBy=multi-user.target
